# cron/crontab
*/1	*	*	*	*	root *applications/testejenkins/cron/background.py
